# barktown-ingest.service
[Unit]
Description=barktown ingest service - watches MinIO /upload-here/ and processes audio files
After=network-online.target minio.service
Wants=network-online.target

[Service]
Type=simple
User=pi

WorkingDirectory=/home/pi/git/enchartme/barktown-ingest
ExecStart=/usr/bin/node ingest-service.mjs
Restart=on-failure
RestartSec=10s

# Edit .env in the WorkingDirectory to set credentials and overrides (see .env.example)
EnvironmentFile=/home/pi/git/enchartme/barktown-ingest/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=barktown-ingest

[Install]
WantedBy=multi-user.target
